# deployment/keycast-signer.service
[Unit]
Description=Keycast NIP-46 Signer Daemon
Documentation=https://github.com/rabble/keycast
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=keycast
Group=keycast
WorkingDirectory=/opt/keycast

# Environment
Environment="RUST_LOG=info"
Environment="DATABASE_PATH=/opt/keycast/database/keycast.db"
Environment="MASTER_KEY_PATH=/opt/keycast/master.key"
Environment="PORT=8080"
# For GCP KMS, add:
#Environment="USE_GCP_KMS=true"
#Environment="GCP_PROJECT_ID=your-project-id"

# Binary
ExecStart=/opt/keycast/keycast_signer

# Restart policy
Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

# Crash handling
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30s

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/keycast/database
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=keycast-signer

[Install]
WantedBy=multi-user.target
